# config/rc.local
#!/usr/bin/bash
export PATH="/bin:/sbin:/usr/local/bin:/usr/local/sbin"

# Source configuration files from /etc/rc.local.d
for file in /etc/rc.local.d/*.sh; do
	if [ -r "$file" ]; then
		bash $file
	fi
done
unset file
sleep 3
grep 0 /sys/kernel/debug/usb/ci_hdrc.0/device | grep speed
if [ $? -eq 0 ]
then
echo host > /sys/kernel/debug/usb/ci_hdrc.0/role
fi
exit 0
